# etc/systemd/system/overviewer.timer
[Unit]
Description=Keep the WHS Minecraft Map up-to-date

[Timer]
OnCalendar=*:0/5
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
